# bin/php/mac/php-versions/php@8.1/homebrew.php@8.1.service
[Unit]
Description=Homebrew generated unit for php@8.1

[Install]
WantedBy=default.target

[Service]
Type=simple
ExecStart=/opt/homebrew/opt/php@8.1/sbin/php-fpm --nodaemonize
Restart=always
WorkingDirectory=/opt/homebrew/var
StandardError=append:/opt/homebrew/var/log/php-fpm.log