# bin/php/mac/php-versions/php@8.2/homebrew.php@8.2.service
[Unit]
Description=Homebrew generated unit for php@8.2

[Install]
WantedBy=default.target

[Service]
Type=simple
ExecStart=/opt/homebrew/opt/php@8.2/sbin/php-fpm --nodaemonize
Restart=always
WorkingDirectory=/opt/homebrew/var
StandardError=append:/opt/homebrew/var/log/php-fpm.log